# Payload/rc.local
#!/bin/bash

##################################################################
#        rc.local Change permissions to some system objects      #
#              Developed by sergio melas 2021-23                 #
##################################################################


#retrive variables for screen and keyboard
KEYBKLIGHT=$( cat /usr/Autorotate/KEYBKLIGHT.conf )

#change permissions
chmod 777 $KEYBKLIGHT


